# cron/fetchprice.cron
*/15 * * * * /usr/bin/php /var/www/halfmyrogers.theholyroger.com/getprice.php